# Pi5 Setup/systemd/mvp_bridge.service
[Unit]
Description=HydraVision MVP Bridge (WebSocket to UDP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/Dev/cross_shore_dev/apps/controller
ExecStart=/usr/bin/python3 /home/admin/Dev/cross_shore_dev/apps/controller/mvp_bridge.py
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

