# Pi5 Setup/autostart/hydravision-kiosk.desktop
[Desktop Entry]
Type=Application
Name=HydraVision Kiosk
Comment=Launch HydraVision UI in Chromium kiosk mode
Exec=/bin/bash -lc "xset s off; xset -dpms; xset s noblank; cog --platform=wl --fullscreen file:///home/admin/Dev/cross_shore_dev/apps/controller/mvp_ui.html"
X-GNOME-Autostart-enabled=true

