# Pi5 Setup/autostart/hydravision-kiosk.desktop
[Desktop Entry]
Type=Application
Name=HydraVision Kiosk
Comment=Launch HydraVision UI in Firefox kiosk mode
Exec=/bin/bash -lc "/usr/local/bin/hydravision-kiosk-launch.sh"
X-GNOME-Autostart-enabled=true

